# classpulse.service
[Unit]
Description=ClassPulse Web Application
After=network.target

[Service]
Type=forking
User=michael
Group=michael
WorkingDirectory=/home/michael/projects/class-pulse
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/home/michael/projects/class-pulse/venv/bin"
ExecStart=/bin/bash /home/michael/projects/class-pulse/start.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=classpulse

[Install]
WantedBy=multi-user.target